# infra/azure/vllm-vm/vllm.service
[Unit]
Description=vLLM Inference Server for Inclusify
Documentation=https://docs.vllm.ai/
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/home/azureuser/inclusify

# Environment configuration
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="HF_HOME=/home/azureuser/.cache/huggingface"
Environment="PATH=/home/azureuser/inclusify/.venv/bin:/usr/local/cuda-12.4/bin:/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/cuda-12.4/lib64"

# vLLM serve command
# - Model: GPTQ-quantized suzume-llama-3-8B
# - dtype: float16 (T4 supports float16, not bfloat16)
# - max-model-len: 2048 (reduced from 4096 to save VRAM)
# - host: 127.0.0.1 (localhost only, no public access)
# - port: 8001 (internal vLLM API)
# - enable-lora: Load LoRA adapters dynamically
# - lora-modules: Inclusify fine-tuned adapter
# - gpu-memory-utilization: 0.9 (90% of VRAM for model + KV cache)
ExecStart=/home/azureuser/inclusify/.venv/bin/python -m vllm.entrypoints.openai.api_server \
    --model /home/azureuser/models/suzume-llama-3-8B-gptq \
    --dtype float16 \
    --max-model-len 2048 \
    --host 127.0.0.1 \
    --port 8001 \
    --enable-lora \
    --lora-modules inclusify=/home/azureuser/inclusify/ml/LoRA_Adapters \
    --gpu-memory-utilization 0.9 \
    --trust-remote-code

# Restart configuration
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/azureuser/models /home/azureuser/.cache
PrivateTmp=true

[Install]
WantedBy=multi-user.target
