# infra/azure/vllm-vm/vllm.service
[Unit]
Description=vLLM Inference Server (Qwen2.5-3B + Inclusify LoRA)
Documentation=https://docs.vllm.ai
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/home/azureuser

# Environment configuration
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="PATH=/home/azureuser/vllm-venv/bin:/usr/local/cuda-12.4/bin:/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/cuda-12.4/lib64"

# vLLM serve command with LoRA adapter (THROUGHPUT OPTIMIZED)
# - FP16 model (5.8GB on disk, ~6GB VRAM)
# - dtype=half (not bfloat16 - T4 incompatible)
# - max-model-len=4096 (doubled from 2048; fits on T4 16GB, needed for chunk-level analysis with full system prompt)
# - host=0.0.0.0 (allow external connections for backend integration)
# - gpu-memory-utilization=0.88 (increased KV cache, headroom for CUDA graphs)
# - max-num-seqs=16 (T4 CUDA graph limit, higher causes OOM during capture)
# - max-num-batched-tokens=8192 (larger token batches for chunked prefill)
# - enable-lora: Load LoRA adapters
# - lora-modules: Inclusify fine-tuned adapter (79.11% accuracy)
# - max-lora-rank=16: Support LoRA rank up to 16 (our adapter uses r=8)
ExecStart=/home/azureuser/vllm-venv/bin/vllm serve \
    /home/azureuser/models/Qwen2.5-3B-Instruct \
    --dtype half \
    --max-model-len 4096 \
    --port 8001 \
    --host 0.0.0.0 \
    --gpu-memory-utilization 0.88 \
    --max-num-seqs 16 \
    --max-num-batched-tokens 8192 \
    --enable-lora \
    --lora-modules inclusify=/home/azureuser/inclusify/ml/adapters/qwen_r8_d0.2 \
    --max-lora-rank 16

# Restart configuration
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm

[Install]
WantedBy=multi-user.target
