# infra/azure/vllm-vm/vllm.service
[Unit]
Description=Inclusify Inference Service (Transformers + LoRA)
Documentation=https://github.com/inclusify
After=network.target

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/home/azureuser/inclusify

# Environment configuration
Environment="CUDA_VISIBLE_DEVICES=0"
Environment="HF_HOME=/home/azureuser/.cache/huggingface"
Environment="PATH=/home/azureuser/inclusify/.venv/bin:/usr/local/cuda-12.4/bin:/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/cuda-12.4/lib64"

# Inference service using transformers directly (T4 compatible)
# - Uses bitsandbytes 4-bit quantization
# - Loads LoRA adapter from ml/LoRA_Adapters
# - Exposes OpenAI-compatible API on port 8001
# - Binds to localhost only (no public access)
ExecStart=/home/azureuser/inclusify/.venv/bin/python /home/azureuser/inclusify/ml/inference_service.py \
    --host 127.0.0.1 \
    --port 8001 \
    --adapter-path /home/azureuser/inclusify/ml/LoRA_Adapters

# Restart configuration
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vllm

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/azureuser/models /home/azureuser/.cache
PrivateTmp=true

[Install]
WantedBy=multi-user.target
